# templates/default/upstart-sentry-celery.conf.erb
# upstart for sentry-celery
#
# Generated by Chef for <%= node[:fqdn] %>
# Local modifications will be overwritten.
#

description "Sentry Celery"

start on started sentry
stop on stopped sentry
#Send KILL after 5 seconds
kill timeout 5
respawn

<%- case node["platform_family"] -%>
<%- when "debian" -%>
setuid <%= @user %>
setgid <%= @group %>
<%- end -%>

env VENV="<%= @work_dir %>"
env LOGDIR="<%= @log_dir %>"
env SENTRY_CMD="$VENV/bin/sentry --config=$VENV/sentry.conf.py celery worker -B -s $LOGDIR/celerybeat-schedule"

pre-start script
  test -f $VENV/bin/sentry || test -f $VENV/sentry.conf.py || { stop ; exit 0; }
end script

script
  sleep 2
  exec >> $LOGDIR/sentry-celery-console.log 2>&1
<%- case node["platform_family"] -%>
<%- when "debian" -%>
  exec $SENTRY_CMD
<%- when "rhel", "fedora", "suse" -%>
  exec sudo -u <%= @user %> -g <%= @group %> sh -c "$SENTRY_CMD"
<%- end -%>
end script
