# templates/default/upstart-sentry-udp.conf.erb
# upstart for sentry udp
#
# Generated by Chef for <%= node[:fqdn] %>
# Local modifications will be overwritten.
#

description "Sentry UDP"

start on started sentry
stop on stopped sentry
#Send KILL after 5 seconds
kill timeout 5
respawn

<%- case node["platform_family"] -%>
<%- when "debian" -%>
setuid <%= @user %>
setgid <%= @group %>
<%- end -%>

env VENV="<%= @work_dir %>"
env SENTRY_CMD="<%= @work_dir %>/bin/sentry --config=<%= @work_dir %>/sentry.conf.py start udp"

pre-start script
  test -f $VENV/bin/sentry || test -f $VENV/sentry.conf.py || { stop ; exit 0; }
end script

script
  sleep 2
  exec >> <%= @log_dir %>/sentry-udp-console.log 2>&1
<%- case node["platform_family"] -%>
<%- when "debian" -%>
  exec $SENTRY_CMD
<%- when "rhel", "fedora", "suse" -%>
  exec su -s /bin/sh -c 'exec "$0" "$@"' <%= @user %> -- $SENTRY_CMD
<%- end -%>
end script
