# templates/default/upstart-rhodecode.conf.erb
# upstart for rhodecode
#
# Generated by Chef for <%= node[:fqdn] %>
# Local modifications will be overwritten.
#

description "Gunicorn for rhodecode"

start on runlevel [2345]
stop on runlevel [!2345]
#Send KILL after 5 seconds
kill timeout 5
respawn

env VENV="<%= @work_dir %>"
env HGENCODING="UTF-8"
export HGENCODING

pre-start script
  test -f $VENV/bin/gunicorn_paster || tset -f $VENV/production.ini || { stop ; exit 0; }
end script

script
  exec $VENV/bin/gunicorn_paster --user <%= @user %> --preload -w <%= @workers %> --log-level <%= @log_level %> --log-file $VENV/gunicorn.log -b <%= @bind %> $VENV/production.ini
end script

